# config/supervisor.gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn staphopia.wsgi --reload -c /staphopia/ebs/staphopia.com/config/gunicorn.conf.py --env DJANGO_SETTINGS_MODULE=staphopia.settings.local
directory=/staphopia/ebs/staphopia.com
user=staphopia
autorestart=true
redirect_stderr=true
